# ds-led.service
[Unit]
Description=DualSense led controller daemon

[Service]
ExecStart=python -m ds_led
Restart=always

[Install]
WantedBy=multi-user.target